# systemd/command-meta-scan.timer
[Unit]
Description=Run Command telemetry anomaly scan hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
