# systemd/synaplex-score.timer
[Unit]
Description=Run synaplex scoring once per intake cadence (intake+8min)
Requires=synaplex-score.service

[Timer]
# Fire 6×/day at :25, aligned to intake's 4h cadence at :17 (so 00:25
# lags 00:17 by ~8 min — plenty of slack for intake's ~70s fetch).
#
# History: this timer originally ran hourly (12×/day) at :20. The :20
# slot fixed the score:stuck-at-midnight issue (intake at 00:17 + score
# at 00:20). But hourly scoring runs 12×/day against a corpus that only
# changes 3-6×/day (intake fires every 4h), so 6-9 of 12 daily score
# runs re-score unchanged data. Free under the heuristic scorer; ~8100
# Sonnet API calls/day under the LLM scorer (when ANTHROPIC_API_KEY
# lands). Aligning to intake cadence eliminates the redundancy without
# losing any scored coverage. See URGENT-synaplex-score-cron-3rd-cycle-
# 2026-05-10.md for the carry-forward audit trail.
OnCalendar=*-*-* 00,04,08,12,16,20:25:00
RandomizedDelaySec=60
Persistent=true
Unit=synaplex-score.service

[Install]
WantedBy=timers.target
