# systemd/workspace-reflect.timer
[Unit]
Description=Run workspace per-project reflection every 12 hours

[Timer]
OnCalendar=*-*-* 02,14:17:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
